# deploy/trading.service
[Unit]
Description=Poly-Maker Trading Bot
After=network.target

[Service]
Type=simple
User=polymaker
WorkingDirectory=/home/polymaker/poly-maker
ExecStart=/home/polymaker/.local/bin/uv run python main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/polymaker/trading.log
StandardError=append:/var/log/polymaker/trading.log

# Environment
Environment=PATH=/home/polymaker/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/polymaker/poly-maker/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/polymaker/poly-maker/positions
PrivateTmp=true

[Install]
WantedBy=multi-user.target
